# conf/systemd.service
[Unit]
Description=Decypharr - Debrid blackhole client
After=network-online.target systemd-udev-settle.service
Wants=network-online.target
Requires=/dev/fuse

[Service]
Type=simple
User=__APP__
Group=__APP__
WorkingDirectory=__INSTALL_DIR__
ExecStart=__INSTALL_DIR__/decypharr --config __DATA_DIR__
Restart=always
RestartSec=3
KillSignal=SIGINT

RuntimeDirectory=__APP__
Environment="HOME=__DATA_DIR__"
Environment="XDG_RUNTIME_DIR=/run/__APP__"

# Required for FUSE mounting (DFS mode)
AmbientCapabilities=CAP_SYS_ADMIN CAP_FSETID CAP_SETGID CAP_SETUID
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_FSETID CAP_SETGID CAP_SETUID
SupplementaryGroups=fuse
DeviceAllow=/dev/fuse rw
MountFlags=shared

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
LockPersonality=yes

StandardOutput=append:/var/log/__APP__/__APP__.log
StandardError=append:/var/log/__APP__/__APP__.log

[Install]
WantedBy=multi-user.target
