# conf/systemd.service
[Unit]
Description=Decypharr - Debrid blackhole client
After=network.target
# If mounting is used, ensure FUSE is available
After=sys-fs-fuse-connections.mount

[Service]
Type=simple
User=__APP__
Group=__APP__

WorkingDirectory=__DATA_DIR__
ExecStart=__INSTALL_DIR__/decypharr --config __DATA_DIR__/config.json

# Restart on failure with backoff
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# FUSE / mount capability (required when mount_type != none)
# These are harmless when mounting is not used.
AmbientCapabilities=CAP_SYS_ADMIN
CapabilityBoundingSet=CAP_SYS_ADMIN

# Allow FUSE mounts to be visible outside the service namespace
MountFlags=shared

# Environment
Environment="HOME=__DATA_DIR__"

# Logging
StandardOutput=append:/var/log/__APP__/__APP__.log
StandardError=append:/var/log/__APP__/__APP__.log

[Install]
WantedBy=multi-user.target
